# conf/supervisor/coplayer-socket.ini
[program:coplayer-socket]
command=node index.js
directory=/home/coplayer-socket
autostart=true
autorestart=true
user=root
startsecs=10
priority=100
stopsignal=QUIT
stdout_logfile=/home/coplayer-socket/logs/coplayer-socket_stdout.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/home/coplayer-socket/logs/coplayer-socket_stderr.log
stderr_logfile_maxbytes=50MB

