# conf/supervisor/lovewangmiao.ini
[program:lovewangmiao_api]
command=/home/venv/bin/uwsgi /etc/uwsgi/lovewangmiao.ini
autostart=true
autorestart=true
user=root
startsecs=10
priority=100
stopsignal=QUIT
stdout_logfile=/home/lovewangmiao/logs/lovewangmiao_api_stdout.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/home/lovewangmiao/logs/lovewangmiao_api_stderr.log
stderr_logfile_maxbytes=50MB
