# conf/supervisor/lovewangmiao.ini
[program:lovewangmiao_api]
command=/home/venv/bin/uwsgi /etc/uwsgi/lovewangmiao.ini
autostart=true
autorestart=true
user=root
startsecs=10
priority=100
stopsignal=QUIT
